# apps/fleet-watcher/fleet-watcher.service
[Unit]
# Fleet watcher — alertas Telegram. Corre en el VPS iagrow-vps (usuario kowen).
Description=Kowen Fleet Watcher (alertas Telegram)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kowen
WorkingDirectory=/home/kowen/proyectos/vending-kowen/apps/fleet-watcher
ExecStart=/usr/bin/python3 /home/kowen/proyectos/vending-kowen/apps/fleet-watcher/watcher.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
